# vessel.service
[Unit]
Description=Vessel - Self-hosted App Deployment Manager
Documentation=https://github.com/Yash121l/Vessel
After=network.target docker.service
Requires=docker.service

[Service]
Type=simple
User=root
ExecStart=/usr/local/bin/vessel serve
Restart=always
RestartSec=3s
StartLimitIntervalSec=0
Environment=VESSEL_CONFIG=/etc/vessel/config.yaml
StandardOutput=journal
StandardError=journal
SyslogIdentifier=vessel

[Install]
WantedBy=multi-user.target
